# bash/profile.d/gpg.sh
#
# ${XDG_CONFIG_HOME}/bash/profile.d/gpg.sh
#

export GNUPGHOME="${XDG_DATA_HOME}/gnupg"
export GPG_TTY="$(tty)"
